# infra/units/etc-secrets.mount
[Mount]
What=/dev/disk/by-id/scsi-0Google_PersistentDisk_persistent-disk-1
Where=/etc/secrets
Type=ext4

[Install]
RequiredBy=local-fs.target
